# net/byedpi-qt6/pkg-plist

bin/byedpi-qt6
bin/ciadpi
etc/rc.d/ciadpi
share/applications/byedpi-qt6.desktop
share/pixmaps/byedpi.png
