# xboxdrv/pkg-plist

bin/xboxdrv
bin/xboxdrvctl
etc/default/xboxdrv
etc/rc.d/xboxdrv
share/man/man1/xboxdrv.1.gz
