# www/opera-proxy/pkg-plist

bin/opera-proxy
bin/operaproxy-qt6
etc/rc.d/operaproxy
share/applications/opera-proxy.desktop
share/pixmaps/opera-proxy.jpg
